# init/twitter-stream.conf
start on runlevel [2345]

respawn
respawn limit 5 60

setuid vagrant
chdir /vagrant/cli
exec php twitter.php 
